# monitor/halo-ddos-monitor.service
[Unit]
Description=Halo CE DDoS Monitor → Discord
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root

# Required: your Discord webhook URL
Environment="DISCORD_WEBHOOK=https://discord.com/api/webhooks/REPLACE_ME/REPLACE_ME"

# Required: one entry per Halo server you run, format "port:Display Name[:max_players]"
# Comma-separated, no quotes around names. Max defaults to 16 if omitted.
# Example with three servers and per-server max:
#   "HALO_SERVERS=2302:Public:8,2303:Scrims:4,2304:Test:16"
Environment="HALO_SERVERS=2302:Server 1:16"

# Optional: comma-separated systemd units to watch for crashes
Environment="WATCHED_SERVICES=halo-server"

# Optional: connect-command helper.
# Public IP is auto-detected via icanhazip.com if PUBLIC_IP is unset.
# SERVER_IPS overrides per port (use when servers live on different hosts).
# SERVER_PASSWORDS embeds the join password into the connect command.
# Environment="PUBLIC_IP=1.2.3.4"
# Environment="SERVER_IPS=2302:1.2.3.4,2303:5.6.7.8"
# Environment="SERVER_PASSWORDS=2303:scrim-pw"

# Optional: proxycheck.io API key (sign up at https://proxycheck.io/dashboard)
# Without a key, free tier is 100 VPN lookups/day. With a key: 1000/day.
# Environment="PROXYCHECK_KEY=your-key-here"

# Optional: tweak detection thresholds (defaults shown)
# Environment="PPS_THRESHOLD=3000"
# Environment="BPS_THRESHOLD=8388608"
# Environment="UNIQUE_IPS_THRESHOLD=40"

Environment="PYTHONUNBUFFERED=1"

ExecStart=/usr/bin/python3 /opt/halo-monitor/halo_ddos_monitor.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
